# geogit.conf
[program:geogit]
directory=/GeoGit/src/parent
user=root
command=mvn jetty:run -pl ../web/app -f pom.xml -Dorg.geogit.web.repository=/GeoGitRepo
autorestart=true
stopsignal=INT
